# deployment/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:octane]
command=php /app/artisan octane:frankenphp --host=0.0.0.0 --port=8000 --workers=4 --max-requests=500
directory=/app
user=www-data
autostart=true
autorestart=true
priority=1
startsecs=5
startretries=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/www-data",USER="www-data"

[program:ssr]
command=php /app/artisan inertia:start-ssr
directory=/app
user=www-data
autostart=true
autorestart=true
priority=2
startsecs=3
startretries=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/www-data",USER="www-data"

[program:queue]
command=php /app/artisan queue:work --sleep=3 --tries=3 --timeout=90 --max-jobs=1000 --max-time=3600
directory=/app
user=www-data
autostart=true
autorestart=true
priority=3
numprocs=2
process_name=%(program_name)s_%(process_num)02d
startsecs=3
stopwaitsecs=90
stopasgroup=true
killasgroup=true
stdout_logfile=/app/storage/logs/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/storage/logs/queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/www-data",USER="www-data"

[program:schedule]
command=/bin/bash -c "while true; do php /app/artisan schedule:run --verbose --no-interaction && sleep 60; done"
directory=/app
user=www-data
autostart=true
autorestart=true
priority=4
startsecs=0
stdout_logfile=/app/storage/logs/schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/storage/logs/schedule-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/home/www-data",USER="www-data"